# Bash/monitoring/ssh/ssh-check.cron
0 21 * * * root /home/vadmin/scripts/monitoring/ssh/ssh.sh -m virtualcollaboard@gmail.com
